# root/etc/services.d/readarr/run
#!/usr/bin/with-contenv bash

umask ${UMASK_SET:-022}

redirect_cmd() {
	if [[ "${DEBUG}" == "true" ]]; then
		"$@"
	else
		"$@" >/dev/null 2>&1
	fi
}

echo "Starting Readarr."
redirect_cmd exec s6-setuidgid abc \
	/app/readarr/bin/Readarr \
	-nobrowser \
	-data=/config
